# .config/autostart/firefox.desktop
[Desktop Entry]
Exec=firefox %u
Type=Application
StartupNotify=false
Terminal=false
