# .config/autostart/firefox.desktop
[Desktop Entry]
Exec=firefox %u
Terminal=false
Type=Application
StartupNotify=false
