# .config/autostart/gnome-terminal.desktop
[Desktop Entry]
Exec=gnome-terminal --maximize
Type=Application
StartupNotify=false
