# .config/autostart/vscode.desktop
[Desktop Entry]
Exec=code %u
Type=Application
StartupNotify=false
Terminal=false
